# snake-game.service
[Unit]
Description=Snake Game Web Service
After=network.target

[Service]
User=your_username
WorkingDirectory=/path/to/snake-game
Environment="PATH=/path/to/your/venv/bin"
ExecStart=gunicorn --workers 1 --bind 0.0.0.0:5001 snake_game:app

[Install]
WantedBy=multi-user.target